# deploy/kyllenios-core.service
[Unit]
Description=kyllenios-core - AI API Gateway Platform
Documentation=https://github.com/ca0fgh/kyllenios-core
After=network.target postgresql.service redis.service
Wants=postgresql.service redis.service

[Service]
Type=simple
User=kyllenios-core
Group=kyllenios-core
WorkingDirectory=/opt/kyllenios-core
ExecStart=/opt/kyllenios-core/kyllenios-core
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=kyllenios-core

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/opt/kyllenios-core

# Environment - Server configuration
# Modify these values to change listen address and port
Environment=GIN_MODE=release
Environment=SERVER_HOST=0.0.0.0
Environment=SERVER_PORT=8080

[Install]
WantedBy=multi-user.target
